# etc/services.d/smbd/run
#!/command/execlineb -P

/usr/sbin/smbd --no-process-group --foreground
