# etc/services.d/nmbd/run
#!/command/execlineb -P

/usr/sbin/nmbd --no-process-group --foreground
